# crontab
0 6 * * * cd COVID-19 && git pull origin master && cd ..
# this will pull the repository at 6 am every day

